# gdrive/.config/systemd/user/gdrive-sync.service
[Unit]
Description=rclone bisync Google Drive

[Service]
Type=oneshot
ExecStart=%h/bin/gdrive-sync
